# config/service/log/run
#!/bin/bash
exec setuidgid jibot sh -c 'exec -a multilog-jibot /usr/local/bin/multilog t ~jibot/jibot/config/service/log/main'
